# shared/launchd/com.nexus.roadmap-sync.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.nexus.roadmap-sync</string>

    <!--
      L1 + L2 enforcement of R-ROADMAP-SYNC.
      L1 WatchPaths: fsevents on roadmap directories → immediate sync.
      L2 StartInterval: 300s safety net in case WatchPaths drops an event.
      H-NOBLOCK: script is idempotent; clean tree → instant no-op.

      Bootstrap (user must invoke manually):
        launchctl bootload gui/$(id -u) \
          /Users/ghost/Dev/nexus/shared/launchd/com.nexus.roadmap-sync.plist
        # or symlink into ~/Library/LaunchAgents and `launchctl load`
    -->

    <key>ProgramArguments</key>
    <array>
        <string>/Users/ghost/Dev/nexus/shared/bin/hexa</string>
        <string>/Users/ghost/Dev/nexus/shared/harness/roadmap_sync.hexa</string>
    </array>

    <!-- L1: fsevents on roadmap dirs (fires within ~1s of any file change) -->
    <key>WatchPaths</key>
    <array>
        <string>/Users/ghost/Dev/nexus/shared/roadmaps</string>
        <string>/Users/ghost/Dev/n6-architecture/shared/roadmaps</string>
    </array>

    <!-- L2: 5-minute safety net -->
    <key>StartInterval</key>
    <integer>300</integer>

    <!-- Batch rapid edits into one sync pass -->
    <key>ThrottleInterval</key>
    <integer>5</integer>

    <key>RunAtLoad</key>
    <false/>

    <key>StandardOutPath</key>
    <string>/Users/ghost/Dev/nexus/shared/logs/roadmap-sync.stdout.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/ghost/Dev/nexus/shared/logs/roadmap-sync.stderr.log</string>

    <key>WorkingDirectory</key>
    <string>/Users/ghost/Dev/nexus</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
        <key>HOME</key>
        <string>/Users/ghost</string>
        <key>LANG</key>
        <string>en_US.UTF-8</string>
        <key>NEXUS</key>
        <string>/Users/ghost/Dev/nexus</string>
    </dict>

    <key>ProcessType</key>
    <string>Background</string>

    <key>Nice</key>
    <integer>10</integer>

    <key>LowPriorityIO</key>
    <true/>
</dict>
</plist>
