# shared/launchd/com.nexus.auto-discovery.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.nexus.auto-discovery</string>

    <!--
      DISC-P3-2 — auto_discovery + growth_registry KPI 자동 업데이트.
      DISABLED 2026-04-19: both referenced python files removed in R37 sweep —
        - auto_discovery_expand.py → commit 3c1ef431 (R37 sweep 2)
        - growth_kpi_update.py → commit 4cac122d (R37 py ban)
      Do NOT bootstrap until hexa-native replacements arrive.
      Original command (archival only, program paths invalid):
        python3 shared/discovery/auto_discovery_expand.py
        && python3 shared/harness/growth_kpi_update.py --once
    -->

    <key>Disabled</key>
    <true/>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>-c</string>
        <string>/usr/bin/python3 /Users/ghost/Dev/nexus/shared/discovery/auto_discovery_expand.py &amp;&amp; /usr/bin/python3 /Users/ghost/Dev/nexus/shared/harness/growth_kpi_update.py --once</string>
    </array>

    <key>StartInterval</key>
    <integer>14400</integer>

    <key>ThrottleInterval</key>
    <integer>600</integer>

    <key>RunAtLoad</key>
    <false/>

    <key>StandardOutPath</key>
    <string>/Users/ghost/Dev/nexus/shared/logs/auto-discovery.stdout.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/ghost/Dev/nexus/shared/logs/auto-discovery.stderr.log</string>

    <key>WorkingDirectory</key>
    <string>/Users/ghost/Dev/nexus</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
        <key>HOME</key>
        <string>/Users/ghost</string>
        <key>LANG</key>
        <string>en_US.UTF-8</string>
        <key>NEXUS</key>
        <string>/Users/ghost/Dev/nexus</string>
    </dict>

    <key>ProcessType</key>
    <string>Background</string>

    <key>Nice</key>
    <integer>15</integer>

    <key>LowPriorityIO</key>
    <true/>
</dict>
</plist>
