# shared/launchd/com.nexus.atlas-drift-monitor.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.nexus.atlas-drift-monitor</string>

    <!--
      ATLAS-P3-3 — hourly atlas drift sampler.
      DISABLED 2026-04-19: referenced atlas_drift_monitor.py was removed in R37
      sweep (commit 3c1ef431, python purge policy). Do NOT bootstrap until
      hexa-native drift monitor replaces it.
      Original command (archival only, program path invalid):
        /opt/homebrew/bin/python3 shared/n6/atlas_drift_monitor.py --sample
    -->

    <key>Disabled</key>
    <true/>

    <key>ProgramArguments</key>
    <array>
        <string>/opt/homebrew/bin/python3</string>
        <string>/Users/ghost/Dev/nexus/shared/n6/atlas_drift_monitor.py</string>
        <string>--sample</string>
    </array>

    <key>StartInterval</key>
    <integer>3600</integer>

    <key>RunAtLoad</key>
    <false/>

    <key>ThrottleInterval</key>
    <integer>60</integer>

    <key>StandardOutPath</key>
    <string>/Users/ghost/Dev/nexus/shared/logs/atlas-drift-monitor.stdout.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/ghost/Dev/nexus/shared/logs/atlas-drift-monitor.stderr.log</string>

    <key>WorkingDirectory</key>
    <string>/Users/ghost/Dev/nexus</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
        <key>HOME</key>
        <string>/Users/ghost</string>
        <key>LANG</key>
        <string>en_US.UTF-8</string>
        <key>NEXUS</key>
        <string>/Users/ghost/Dev/nexus</string>
    </dict>

    <key>ProcessType</key>
    <string>Background</string>

    <key>Nice</key>
    <integer>15</integer>

    <key>LowPriorityIO</key>
    <true/>
</dict>
</plist>
