# shared/launchd/com.nexus.auto-discovery.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.nexus.auto-discovery</string>

    <!--
      DISC-P3-2 — 영구 백그라운드 auto_discovery + growth_registry KPI 자동 업데이트.
      14400s (4h) 주기로 run_auto_discovery.sh 실행:
        1) python3 shared/discovery/auto_discovery_expand.py         (100-cycle 확장 + breakthrough merge)
        2) python3 shared/harness/growth_kpi_update.py --once        (5 KPI 갱신 원자 쓰기)
      NOTE: DO NOT launchctl load in-repo. 수동 활성화는 운영자 책임 —
            `launchctl load ~/Library/LaunchAgents/com.nexus.auto-discovery.plist` 요함.
    -->

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>-c</string>
        <string>/usr/bin/python3 /Users/ghost/Dev/nexus/shared/discovery/auto_discovery_expand.py &amp;&amp; /usr/bin/python3 /Users/ghost/Dev/nexus/shared/harness/growth_kpi_update.py --once</string>
    </array>

    <key>StartInterval</key>
    <integer>14400</integer>

    <key>ThrottleInterval</key>
    <integer>600</integer>

    <key>RunAtLoad</key>
    <false/>

    <key>StandardOutPath</key>
    <string>/Users/ghost/Dev/nexus/shared/logs/auto-discovery.stdout.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/ghost/Dev/nexus/shared/logs/auto-discovery.stderr.log</string>

    <key>WorkingDirectory</key>
    <string>/Users/ghost/Dev/nexus</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
        <key>HOME</key>
        <string>/Users/ghost</string>
        <key>LANG</key>
        <string>en_US.UTF-8</string>
        <key>NEXUS</key>
        <string>/Users/ghost/Dev/nexus</string>
    </dict>

    <key>ProcessType</key>
    <string>Background</string>

    <key>Nice</key>
    <integer>15</integer>

    <key>LowPriorityIO</key>
    <true/>
</dict>
</plist>
