# shared/launchd/com.nexus.mac-defense.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.nexus.mac-defense</string>

    <!--
      Mac kernel panic 방어 — apfsd watchdog timeout 재발 차단.

      검증 후 활성화:
        1) ./mac_defense_watchdog 수동 1회 실행 → ~/Library/Logs/mac_defense.log 확인
        2) Disabled=false 로 변경
        3) cp ~/Dev/nexus/shared/launchd/com.nexus.mac-defense.plist ~/Library/LaunchAgents/
        4) launchctl bootstrap gui/$UID ~/Library/LaunchAgents/com.nexus.mac-defense.plist
        5) 24h dry-run log 관찰 후 MAC_DEFENSE_KILL=1 활성화 결정

      cf. 2026-04-20 03:33 panic incident — apfsd 77% CPU/117s → watchdog 91s timeout.
    -->

    <key>Disabled</key>
    <true/>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Users/ghost/Dev/nexus/shared/scripts/bin/mac_defense_watchdog</string>
    </array>

    <key>StartInterval</key>
    <integer>10</integer>

    <key>RunAtLoad</key>
    <false/>

    <key>ThrottleInterval</key>
    <integer>10</integer>

    <key>StandardOutPath</key>
    <string>/Users/ghost/Library/Logs/mac_defense.stdout.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/ghost/Library/Logs/mac_defense.stderr.log</string>

    <key>WorkingDirectory</key>
    <string>/Users/ghost</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
        <key>HOME</key>
        <string>/Users/ghost</string>
        <key>LANG</key>
        <string>en_US.UTF-8</string>
        <key>MAC_DEFENSE_KILL</key>
        <string>0</string>
    </dict>

    <key>ProcessType</key>
    <string>Background</string>

    <key>Nice</key>
    <integer>10</integer>

    <key>LowPriorityIO</key>
    <true/>
</dict>
</plist>
