# shared/launchd/com.nexus.sensor.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.nexus.sensor</string>

    <key>ProgramArguments</key>
    <array>
        <string>/Users/ghost/Dev/nexus/shared/bin/hexa</string>
        <string>run</string>
        <string>/Users/ghost/Dev/nexus/shared/harness/sensor.hexa</string>
        <string>all</string>
    </array>

    <key>StartInterval</key>
    <integer>3600</integer>

    <key>ThrottleInterval</key>
    <integer>600</integer>

    <key>RunAtLoad</key>
    <false/>

    <key>StandardOutPath</key>
    <string>/Users/ghost/Dev/nexus/shared/logs/sensor.stdout.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/ghost/Dev/nexus/shared/logs/sensor.stderr.log</string>

    <key>WorkingDirectory</key>
    <string>/Users/ghost/Dev/nexus</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
        <key>HOME</key>
        <string>/Users/ghost</string>
        <key>LANG</key>
        <string>en_US.UTF-8</string>
        <key>NEXUS</key>
        <string>/Users/ghost/Dev/nexus</string>
    </dict>

    <key>ProcessType</key>
    <string>Background</string>

    <key>Nice</key>
    <integer>15</integer>

    <key>LowPriorityIO</key>
    <true/>
</dict>
</plist>
