# shared/launchd/com.anima.r2-sync-watcher.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
  com.anima.r2-sync-watcher — host-side R2 sync loop (2h interval)

  Purpose: immune to RUNPOD_VOLUMELESS_EPHEMERAL_WIPE. Pulls ckpt dirs
           from RunPod pods to R2 on a fixed cadence, driven by the
           host (not the pod). Pod restarts cannot kill this schedule.

  Install (macOS host):
    cp $ANIMA/shared/launchd/com.anima.r2-sync-watcher.plist \
       ~/Library/LaunchAgents/
    launchctl load -w ~/Library/LaunchAgents/com.anima.r2-sync-watcher.plist

  Verify:
    launchctl list | grep com.anima.r2-sync-watcher
    tail -f $ANIMA/shared/logs/r2_sync_watcher.log

  Reload after edit:
    launchctl unload ~/Library/LaunchAgents/com.anima.r2-sync-watcher.plist
    launchctl load   ~/Library/LaunchAgents/com.anima.r2-sync-watcher.plist

  Linux host alternative (no launchd):
    crontab -e
    0 */2 * * * /Users/ghost/Dev/nexus/shared/bin/hexa \
      /Users/ghost/Dev/anima/scripts/r2_sync_watcher.hexa \
      >> /Users/ghost/Dev/anima/shared/logs/r2_sync_watcher.log 2>&1

  Config: $ANIMA/shared/config/r2_sync_watchers.json
-->
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.anima.r2-sync-watcher</string>

    <key>ProgramArguments</key>
    <array>
        <string>/Users/ghost/Dev/nexus/shared/bin/hexa</string>
        <string>/Users/ghost/Dev/anima/scripts/r2_sync_watcher.hexa</string>
    </array>

    <key>StartInterval</key>
    <integer>7200</integer>

    <key>ThrottleInterval</key>
    <integer>7200</integer>

    <key>RunAtLoad</key>
    <false/>

    <key>StandardOutPath</key>
    <string>/Users/ghost/Dev/anima/shared/logs/r2_sync_watcher.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/ghost/Dev/anima/shared/logs/r2_sync_watcher.log</string>

    <key>WorkingDirectory</key>
    <string>/Users/ghost/Dev/anima</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/Users/ghost/.local/bin:/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
        <key>HOME</key>
        <string>/Users/ghost</string>
        <key>LANG</key>
        <string>en_US.UTF-8</string>
        <key>NEXUS</key>
        <string>/Users/ghost/Dev/nexus</string>
        <key>ANIMA</key>
        <string>/Users/ghost/Dev/anima</string>
    </dict>

    <key>ProcessType</key>
    <string>Background</string>

    <key>Nice</key>
    <integer>10</integer>

    <key>LowPriorityIO</key>
    <true/>
</dict>
</plist>
