# shared/launchd/com.nexus.discovery-rate.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.nexus.discovery-rate</string>

    <!--
      Mk.X axis 3 — discovery rate sampler (atlas 포화 자동 감지).
      Runs every 30min (G4 cadence, 2026-04-19). Writes shared/n6/discovery_rate.json + appends discovery_rate.jsonl history (14d retention).
      Plateau detection (vertical/horizontal/isolation/critical) + suggested_axis 출력.
      DO NOT launchctl load unless user explicitly requests:
        launchctl bootstrap gui/$(id -u) shared/launchd/com.nexus.discovery-rate.plist
      Validate:
        plutil -lint shared/launchd/com.nexus.discovery-rate.plist
      Manual:
        hexa shared/n6/discovery_rate.hexa
    -->

    <key>ProgramArguments</key>
    <array>
        <string>/Users/ghost/Dev/nexus/shared/bin/hexa</string>
        <string>/Users/ghost/Dev/nexus/shared/n6/discovery_rate.hexa</string>
        <string>--quiet</string>
    </array>

    <key>StartInterval</key>
    <integer>1800</integer>

    <key>RunAtLoad</key>
    <false/>

    <key>ThrottleInterval</key>
    <integer>60</integer>

    <key>StandardOutPath</key>
    <string>/Users/ghost/Dev/nexus/shared/logs/discovery-rate.stdout.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/ghost/Dev/nexus/shared/logs/discovery-rate.stderr.log</string>

    <key>WorkingDirectory</key>
    <string>/Users/ghost/Dev/nexus</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
        <key>HOME</key>
        <string>/Users/ghost</string>
        <key>LANG</key>
        <string>en_US.UTF-8</string>
        <key>NEXUS</key>
        <string>/Users/ghost/Dev/nexus</string>
    </dict>

    <key>ProcessType</key>
    <string>Background</string>

    <key>Nice</key>
    <integer>15</integer>

    <key>LowPriorityIO</key>
    <true/>
</dict>
</plist>
