# config/launchd/com.nexus.psi_calibration_refresh.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
  com.nexus.psi_calibration_refresh — ω-network-6b daily PSI calibration cache refresh.

  Pairs with /Users/ghost/core/nexus/tool/psi_calibration_cache.hexa (ω-network-6).
  ω-network-6 shipped a 24h-TTL cache that only refreshes when an operator
  sources the calibration tool. For unattended ops (cron-driven hexa_remote
  calls), the cache could remain stale indefinitely. This plist forces a
  refresh once per day so state/psi_calibration_active.tsv stays recent
  regardless of operator interaction.

  Cadence: daily at 04:20 (StartCalendarInterval Hour=4 Minute=20).
    - Off-peak; differs from com.nexus.bridge_pipeline (:15) to avoid
      simultaneous I/O contention.
    - One fire/day; underlying tool is <1s wall.

  Manual install (NOT performed by ω-network-6b — operator-initiated only):
    cp /Users/ghost/core/nexus/config/launchd/com.nexus.psi_calibration_refresh.plist \
       ~/Library/LaunchAgents/
    mkdir -p ~/.hx/log
    launchctl bootstrap gui/$UID ~/Library/LaunchAgents/com.nexus.psi_calibration_refresh.plist
    launchctl enable     gui/$UID/com.nexus.psi_calibration_refresh
    # next fire is the upcoming 04:20. Force first run if needed:
    launchctl kickstart -k gui/$UID/com.nexus.psi_calibration_refresh

  Discoverable install help:
    hexa run /Users/ghost/core/nexus/tool/psi_calibration_install_help.hexa

  Status / logs:
    launchctl print gui/$UID/com.nexus.psi_calibration_refresh
    tail -F ~/.hx/log/psi_calibration_refresh.log ~/.hx/log/psi_calibration_refresh.err
    tail -n 5 /Users/ghost/core/nexus/state/psi_calibration_active.tsv

  Uninstall:
    launchctl bootout gui/$UID/com.nexus.psi_calibration_refresh
    rm ~/Library/LaunchAgents/com.nexus.psi_calibration_refresh.plist

  Notes:
    - RunAtLoad=false: don't fire on bootstrap; wait for next 04:20.
    - KeepAlive=false: one-shot per schedule, not a daemon.
    - Nice=5: lower priority than bridge_pipeline (Nice=1); housekeeping.
    - LowPriorityIO=true: yield I/O to foreground operator commands.
    - psi_calibration_cache.hexa parses argv positionally via _has_flag();
      --refresh-force is consumed; no HEXA_ARGV plumbing needed.
    - PATH includes /opt/homebrew/bin so the tool's stat/date/sed/tr
      sub-execs resolve consistently with operator shells.
-->
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.nexus.psi_calibration_refresh</string>

  <key>ProgramArguments</key>
  <array>
    <string>/Users/ghost/core/hexa-lang/hexa</string>
    <string>run</string>
    <string>/Users/ghost/core/nexus/tool/psi_calibration_cache.hexa</string>
    <string>--refresh-force</string>
  </array>

  <key>StartCalendarInterval</key>
  <dict>
    <key>Hour</key>
    <integer>4</integer>
    <key>Minute</key>
    <integer>20</integer>
  </dict>

  <key>RunAtLoad</key>
  <false/>

  <key>KeepAlive</key>
  <false/>

  <key>ProcessType</key>
  <string>Background</string>

  <key>Nice</key>
  <integer>5</integer>

  <key>LowPriorityIO</key>
  <true/>

  <key>StandardOutPath</key>
  <string>/Users/ghost/.hx/log/psi_calibration_refresh.log</string>

  <key>StandardErrorPath</key>
  <string>/Users/ghost/.hx/log/psi_calibration_refresh.err</string>

  <key>WorkingDirectory</key>
  <string>/Users/ghost/core/nexus</string>

  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key>
    <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
    <key>HOME</key>
    <string>/Users/ghost</string>
    <!-- HEXA_LANG mirrors com.nexus.bridge_pipeline for resolver fallback;
         psi_calibration_cache.hexa resolves NEXUS_ROOT via `git -C` and
         HEXA_BIN via `command -v hexa` (PATH above) so no further env
         duplication required. -->
    <key>HEXA_LANG</key>
    <string>/Users/ghost/core/hexa-lang</string>
  </dict>
</dict>
</plist>
