# systemd/nexus-launcher-cap.service
[Unit]
Description=AG4 launcher_cap.hexa — pre-cap enforcement watchdog
Documentation=file:///root/airgenome/launcher_cap.hexa
After=network.target

[Service]
Type=simple
ExecStart=/root/Dev/hexa-lang/target/release/hexa /root/airgenome/launcher_cap.hexa
WorkingDirectory=/root/airgenome
Restart=on-failure
RestartSec=10
WatchdogSec=120

# Logging — launcher_cap writes JSON-line to AIRGENOME_CAP_LOG;
# stdout/stderr go to journald by default
Environment=AIRGENOME_CAP_LOG=/var/log/airgenome-cap.log
StandardOutput=journal
StandardError=journal
SyslogIdentifier=launcher-cap

# Hardening
NoNewPrivileges=true
ProtectSystem=strict
ReadWritePaths=/var/log
ProtectHome=no
PrivateTmp=true

[Install]
WantedBy=multi-user.target
