# tool/com.nexus.atlas-meta-scan.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
  com.nexus.atlas-meta-scan — R18 Phase-1 MVP periodic atlas meta-health sampler.
  Pairs with /Users/ghost/core/nexus/tool/atlas_meta_scan.hexa (nxs-20260422-002).

  Cadence: every 12h. Appends one row to state/atlas_health_timeline.jsonl.
  RunAtLoad=false: first run is manual (dry-run verification, see below).

  Manual install (after review):
    cp /Users/ghost/core/nexus/tool/com.nexus.atlas-meta-scan.plist \
       ~/Library/LaunchAgents/
    launchctl bootstrap gui/$UID ~/Library/LaunchAgents/com.nexus.atlas-meta-scan.plist
    launchctl enable     gui/$UID/com.nexus.atlas-meta-scan
    # first-run manual verification (dry):
    /opt/homebrew/bin/hexa /Users/ghost/core/nexus/tool/atlas_meta_scan.hexa --dry-run
    # then kickstart the agent:
    launchctl kickstart -k gui/$UID/com.nexus.atlas-meta-scan

  Status / logs:
    launchctl print gui/$UID/com.nexus.atlas-meta-scan
    tail -F /tmp/nexus_atlas_meta_scan.out.log /tmp/nexus_atlas_meta_scan.err.log

  Uninstall:
    launchctl bootout gui/$UID/com.nexus.atlas-meta-scan
    rm ~/Library/LaunchAgents/com.nexus.atlas-meta-scan.plist

  Notes:
    - hexa binary path: /opt/homebrew/bin/hexa (Apple Silicon Homebrew).
      Override PATH below if you keep hexa elsewhere (/usr/local/bin on Intel,
      ~/core/hexa-lang/hexa for dev builds). ProgramArguments[0] is the
      canonical /usr/local/bin/hexa symlink — if absent on your box,
      replace with the real path before installing.
    - ProcessType=Background → jetsam priority 40 (killed before UI under pressure).
    - KeepAlive=false: this is a periodic sampler, not a daemon. launchd
      re-spawns every StartInterval seconds regardless.
    - RunAtLoad=false: avoid side-effects at install time. First row lands
      12h after bootstrap, or immediately via `launchctl kickstart`.
    - Nice=10: never starves interactive work.
-->
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.nexus.atlas-meta-scan</string>

  <key>ProgramArguments</key>
  <array>
    <string>/usr/local/bin/hexa</string>
    <string>/Users/ghost/core/nexus/tool/atlas_meta_scan.hexa</string>
  </array>

  <key>StartInterval</key>
  <integer>43200</integer>

  <key>RunAtLoad</key>
  <false/>

  <key>KeepAlive</key>
  <false/>

  <key>ProcessType</key>
  <string>Background</string>

  <key>Nice</key>
  <integer>10</integer>

  <key>LowPriorityIO</key>
  <true/>

  <key>StandardOutPath</key>
  <string>/tmp/nexus_atlas_meta_scan.out.log</string>

  <key>StandardErrorPath</key>
  <string>/tmp/nexus_atlas_meta_scan.err.log</string>

  <key>WorkingDirectory</key>
  <string>/Users/ghost/core/nexus</string>

  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key>
    <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
    <key>HOME</key>
    <string>/Users/ghost</string>
  </dict>
</dict>
</plist>
