# tool/com.nexus.beyond-omega-daily.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
  com.nexus.beyond-omega-daily — nxs-20260425-004 cycle 10 daily ghost-ceiling timeline.
  Pairs with /Users/ghost/core/nexus/tool/beyond_omega_ghost_trace.py --cron mode (cycle 5 v4).

  Cadence: daily 03:13 local.
    - --cron = --append + summary 시계열 snapshot
      (state/ghost_ceiling_summary.daily.YYYY-MM-DD.json)
    - default protected mode (NEXUS_BACK_ACTION_ON unset):
      cycle 5 SELF_OUTPUTS skip 활성 → first-order distribution timeline 만 누적.
      cycle 8 second-order echo 는 본 plist 에서 측정하지 않음 (사용자 수동 override only).
    - 기존 trace.jsonl 에 file:lineno dedup 으로 incremental 누적.
    - daily snapshot 으로 7-30 day dataset 구축 → cycle 11+ ordinal hierarchy 정량화 기반.

  Manual install (after review — user 가 승인 후 수동 실행):
    cp /Users/ghost/core/nexus/tool/com.nexus.beyond-omega-daily.plist \
       ~/Library/LaunchAgents/
    launchctl bootstrap gui/$UID ~/Library/LaunchAgents/com.nexus.beyond-omega-daily.plist
    launchctl enable     gui/$UID/com.nexus.beyond-omega-daily
    # first-run manual verification (optional dry):
    /opt/homebrew/bin/python3 /Users/ghost/core/nexus/tool/beyond_omega_ghost_trace.py --cron
    # then kickstart the agent (skip 다음 03:13 까지 기다리지 않으려면):
    launchctl kickstart -k gui/$UID/com.nexus.beyond-omega-daily

  Status / logs:
    launchctl print gui/$UID/com.nexus.beyond-omega-daily
    tail -F /tmp/nexus_beyond_omega_daily.out.log /tmp/nexus_beyond_omega_daily.err.log
    ls -la /Users/ghost/core/nexus/state/ghost_ceiling_summary.daily.*.json

  Uninstall:
    launchctl bootout gui/$UID/com.nexus.beyond-omega-daily
    rm ~/Library/LaunchAgents/com.nexus.beyond-omega-daily.plist

  Notes:
    - StartCalendarInterval (Hour=3, Minute=13) instead of StartInterval — daily fixed time,
      offset from atlas-meta-scan/omega-metrics 12h cadence to avoid I/O overlap.
    - RunAtLoad=false: avoid surprise on install. First snapshot lands next 03:13 or via kickstart.
    - python3 path: /opt/homebrew/bin/python3 (Apple Silicon Homebrew, Python 3.14).
    - probe v4 idempotent (cycle 5 fix) — 안전하게 재실행해도 trace.jsonl 무한 누적 없음.
    - default protected mode 라서 cycle 8 echo (NEXUS_BACK_ACTION_ON=1) 는 trigger 되지 않음.
      L_{ω+1} second-order distribution 시계열을 보려면 별도 plist 또는 수동 발사 필요.
    - Expected accumulation: ~1 daily snapshot/day, summary file size <10KB/snapshot,
      30 days = ~300KB total — gitignore 권고 (ghost_ceiling_summary.daily.*.json).
-->
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.nexus.beyond-omega-daily</string>

  <key>ProgramArguments</key>
  <array>
    <string>/opt/homebrew/bin/python3</string>
    <string>/Users/ghost/core/nexus/tool/beyond_omega_ghost_trace.py</string>
    <string>--cron</string>
  </array>

  <key>StartCalendarInterval</key>
  <dict>
    <key>Hour</key>
    <integer>3</integer>
    <key>Minute</key>
    <integer>13</integer>
  </dict>

  <key>RunAtLoad</key>
  <false/>

  <key>KeepAlive</key>
  <false/>

  <key>ProcessType</key>
  <string>Background</string>

  <key>Nice</key>
  <integer>10</integer>

  <key>LowPriorityIO</key>
  <true/>

  <key>StandardOutPath</key>
  <string>/tmp/nexus_beyond_omega_daily.out.log</string>

  <key>StandardErrorPath</key>
  <string>/tmp/nexus_beyond_omega_daily.err.log</string>

  <key>WorkingDirectory</key>
  <string>/Users/ghost/core/nexus</string>

  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key>
    <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
    <key>HOME</key>
    <string>/Users/ghost</string>
    <!-- default protected mode: cycle 5 SELF_OUTPUTS skip 활성, first-order distribution 만 누적.
         cycle 8 echo 활성화 원하면 NEXUS_BACK_ACTION_ON=1 추가 (별도 plist 권고). -->
  </dict>
</dict>
</plist>
</content>
</invoke>