# tool/com.airgenome.hexa_remote_sync.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
  com.airgenome.hexa_remote_sync — fswatch daemon mirroring airgenome+nexus to ubu.
  Pairs with /Users/ghost/core/nexus/shared/harness/hexa_remote_sync.sh.

  Manual install (after review):
    cp /Users/ghost/core/nexus/shared/harness/com.airgenome.hexa_remote_sync.plist \
       ~/Library/LaunchAgents/
    launchctl bootstrap gui/$UID ~/Library/LaunchAgents/com.airgenome.hexa_remote_sync.plist
    launchctl enable     gui/$UID/com.airgenome.hexa_remote_sync
    launchctl kickstart -k gui/$UID/com.airgenome.hexa_remote_sync

  Status / logs:
    launchctl print gui/$UID/com.airgenome.hexa_remote_sync
    tail -F /tmp/hexa_remote_sync.out.log /tmp/hexa_remote_sync.err.log

  Uninstall:
    launchctl bootout gui/$UID/com.airgenome.hexa_remote_sync
    rm ~/Library/LaunchAgents/com.airgenome.hexa_remote_sync.plist

  Notes:
    - ProcessType=Background gives jetsam priority 40 (vs Terminal's 180),
      so under memory pressure this daemon is killed long before the UI.
      Mac-death protection per safe_hexa_launchd.sh pattern.
    - KeepAlive=true: if fswatch dies (e.g. dir unmounted), launchd respawns.
    - Nice=10: never starves interactive work.
    - PATH includes /opt/homebrew/bin so `fswatch` and `rsync` resolve.
-->
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.airgenome.hexa_remote_sync</string>

  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>/Users/ghost/core/nexus/shared/harness/hexa_remote_sync.sh</string>
  </array>

  <key>RunAtLoad</key>
  <true/>

  <key>KeepAlive</key>
  <true/>

  <key>ProcessType</key>
  <string>Background</string>

  <key>Nice</key>
  <integer>10</integer>

  <key>LowPriorityIO</key>
  <true/>

  <key>ThrottleInterval</key>
  <integer>10</integer>

  <key>StandardOutPath</key>
  <string>/tmp/hexa_remote_sync.out.log</string>

  <key>StandardErrorPath</key>
  <string>/tmp/hexa_remote_sync.err.log</string>

  <key>WorkingDirectory</key>
  <string>/Users/ghost</string>

  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key>
    <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
    <key>HOME</key>
    <string>/Users/ghost</string>
    <key>HEXA_REMOTE_HOST</key>
    <string>ubu</string>
  </dict>
</dict>
</plist>
