# tool/com.nexus.dormancy-wake.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
  com.nexus.dormancy-wake — LaunchAgent polling meta_trigger_watchlist for
  elapsed-threshold triggers (e.g. blowup_dormancy_wake ≥72h).

  Role:
    - Cheap hourly poll (StartInterval=3600). Inner logic (dormancy_wake_tick.hexa)
      re-reads state/meta_trigger_watchlist.jsonl and only fires when the
      row's `fires_when: elapsed>=Nh` threshold has actually elapsed since
      `last_event_ts`. Misses at most 1h; no load cost beyond an hourly tick.
    - Unblocks `no_daemon_hook_infra` (B1_wake / B2 / B4 axes).

  Manual install (after review — DO NOT AUTO-INSTALL):
    cp /Users/ghost/core/nexus/tool/com.nexus.dormancy-wake.plist \
       ~/Library/LaunchAgents/
    launchctl bootstrap gui/$UID ~/Library/LaunchAgents/com.nexus.dormancy-wake.plist
    launchctl enable     gui/$UID/com.nexus.dormancy-wake
    launchctl kickstart -k gui/$UID/com.nexus.dormancy-wake

  Status / logs:
    launchctl print gui/$UID/com.nexus.dormancy-wake
    tail -F /tmp/nexus_dormancy_wake.out.log /tmp/nexus_dormancy_wake.err.log

  Uninstall:
    launchctl bootout gui/$UID/com.nexus.dormancy-wake
    rm ~/Library/LaunchAgents/com.nexus.dormancy-wake.plist

  Notes:
    - hexa binary path hardcoded: /Users/ghost/core/hexa-lang/hexa
      (verified via `which hexa` at author-time).
      /usr/local/bin/hexa does NOT exist on this host; using real path.
    - ProcessType=Background → jetsam priority 40 (killed before UI under memory pressure).
    - KeepAlive=false: this is a 1-shot tick, re-scheduled by StartInterval.
    - Nice=10: never starves interactive work.
    - Inner script is singleton-guarded via /tmp/nexus_dormancy_wake.pid (same
      SHOT_PID_FILE pattern as nexus_growth_daemon.hexa).
-->
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.nexus.dormancy-wake</string>

  <key>ProgramArguments</key>
  <array>
    <string>/Users/ghost/core/hexa-lang/hexa</string>
    <string>/Users/ghost/core/nexus/tool/dormancy_wake_tick.hexa</string>
  </array>

  <key>StartInterval</key>
  <integer>3600</integer>

  <key>RunAtLoad</key>
  <true/>

  <key>KeepAlive</key>
  <false/>

  <key>ProcessType</key>
  <string>Background</string>

  <key>Nice</key>
  <integer>10</integer>

  <key>LowPriorityIO</key>
  <true/>

  <key>StandardOutPath</key>
  <string>/tmp/nexus_dormancy_wake.out.log</string>

  <key>StandardErrorPath</key>
  <string>/tmp/nexus_dormancy_wake.err.log</string>

  <key>WorkingDirectory</key>
  <string>/Users/ghost/core/nexus</string>

  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key>
    <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
    <key>HOME</key>
    <string>/Users/ghost</string>
    <key>NEXUS</key>
    <string>/Users/ghost/core/nexus</string>
  </dict>
</dict>
</plist>
