# tool/dev.hexa-lang.raw81-omega-cycle-monitor.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
  raw 81 omega cycle falsifier monitor — daily cadence + login auto-run.

  Watches the raw 71 falsifier conditions for the hive_test_latency
  omega cycle (state/design_strategy_trawl/2026-04-25_hive_test_latency_*):
    - cache_hit_ratio ≥ 0.95 (floor)
    - fdd_p50 ≤ 1500 ms     (ceiling)

  Violation → ai-native error trailer (raw 66) + audit ledger (raw 77) +
  /tmp/hexa_error_ainative/raw81-monitor-<ts>.jsonl.

  Runs at user login (RunAtLoad) and every 86400 s (24 h) thereafter.

  SSOT: nexus (ω-cycle 영역). Source plist + hexa script live under
  /Users/ghost/core/nexus/tool/. Audit ledger stays under hive (the data
  this tool watches is hive_test_latency log; ledger stays close to the
  measured artifact). 2026-04-26 이관: hive → nexus.

  Install:
    cp /Users/ghost/core/nexus/tool/dev.hexa-lang.raw81-omega-cycle-monitor.plist \
       ~/Library/LaunchAgents/
    launchctl bootstrap gui/$(id -u) ~/Library/LaunchAgents/dev.hexa-lang.raw81-omega-cycle-monitor.plist
    launchctl enable    gui/$(id -u)/dev.hexa-lang.raw81-omega-cycle-monitor

  Reload after path migration (existing install):
    launchctl bootout gui/$(id -u)/dev.hexa-lang.raw81-omega-cycle-monitor
    cp /Users/ghost/core/nexus/tool/dev.hexa-lang.raw81-omega-cycle-monitor.plist \
       ~/Library/LaunchAgents/
    launchctl bootstrap gui/$(id -u) ~/Library/LaunchAgents/dev.hexa-lang.raw81-omega-cycle-monitor.plist

  Logs:
    ~/.hx/log/raw81-omega-monitor.log
    ~/.hx/log/raw81-omega-monitor.err

  Audit ledger (relative to WorkingDirectory):
    ~/core/hive/state/raw81/falsifier_audit.jsonl   (raw 77 schema)

  ai-native error log (on violation only):
    /tmp/hexa_error_ainative/raw81-monitor-<ts>.jsonl   (raw 66)
-->
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>dev.hexa-lang.raw81-omega-cycle-monitor</string>
  <key>ProgramArguments</key>
  <array>
    <string>/Users/ghost/.hx/bin/nexus</string>
    <string>omega-monitor</string>
    <string>check</string>
  </array>
  <key>WorkingDirectory</key>
  <string>/Users/ghost/core/hive</string>
  <key>StartInterval</key>
  <integer>86400</integer>
  <key>RunAtLoad</key>
  <true/>
  <key>StandardOutPath</key>
  <string>/Users/ghost/.hx/log/raw81-omega-monitor.log</string>
  <key>StandardErrorPath</key>
  <string>/Users/ghost/.hx/log/raw81-omega-monitor.err</string>
  <key>EnvironmentVariables</key>
  <dict>
    <key>WS</key>
    <string>/Users/ghost/core</string>
    <key>HEXA_LANG</key>
    <string>/Users/ghost/core/hexa-lang</string>
    <!-- nexus CLI 의 nexus_root() 가 default $HOME/core/nexus 를 보지 않도록
         실제 SSOT 위치를 명시 (run.hexa#nexus_root). -->
    <key>NEXUS</key>
    <string>/Users/ghost/core/nexus</string>
    <key>LC_ALL</key>
    <string>C</string>
    <key>PATH</key>
    <string>/usr/local/bin:/opt/homebrew/bin:/usr/bin:/bin</string>
    <!-- Force Mac-local: tool reads ~/.hive/agent/logs on host filesystem;
         container symlink semantics differ. (Same rationale as raw 45 /
         raw 81-default-model-lock plists.) -->
    <key>HEXA_RESOLVER_NO_REROUTE</key>
    <string>1</string>
  </dict>
  <key>ProcessType</key>
  <string>Background</string>
  <key>LowPriorityIO</key>
  <true/>
  <key>Nice</key>
  <integer>15</integer>
</dict>
</plist>
