# state/ops/sshd-selfheal.service
[Unit]
Description=sshd liveness probe + auto-restart (ubu1 post-2026-04-24 hang defense)
Documentation=file:/usr/local/sbin/sshd_selfheal
After=ssh.service

[Service]
Type=oneshot
ExecStart=/usr/local/sbin/sshd_selfheal
Nice=10
IOSchedulingClass=idle
# Run with minimal privileges but allow systemctl + reboot.
User=root
# Don't fail the timer unit if probe is slow.
TimeoutStartSec=15s
