# launchd/com.nexus.atlas-sync.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.nexus.atlas-sync</string>

    <!-- 15min 주기 diff (pull 은 explicit, alert 만) -->
    <key>ProgramArguments</key>
    <array>
        <string>/Users/ghost/Dev/nexus/shared/bin/hexa</string>
        <string>/Users/ghost/Dev/nexus/shared/scripts/atlas_sync_remote.hexa</string>
        <string>--direction</string>
        <string>diff</string>
        <string>--host</string>
        <string>all</string>
    </array>

    <key>StartInterval</key>
    <integer>900</integer>

    <key>ThrottleInterval</key>
    <integer>900</integer>

    <key>RunAtLoad</key>
    <false/>

    <key>StandardOutPath</key>
    <string>/Users/ghost/Dev/nexus/shared/logs/atlas-sync.stdout.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/ghost/Dev/nexus/shared/logs/atlas-sync.stderr.log</string>

    <key>WorkingDirectory</key>
    <string>/Users/ghost/Dev/nexus</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/Users/ghost/.local/bin:/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
        <key>HOME</key>
        <string>/Users/ghost</string>
        <key>LANG</key>
        <string>en_US.UTF-8</string>
        <key>NEXUS</key>
        <string>/Users/ghost/Dev/nexus</string>
        <!-- ssh credentials 가 mac 에만 — 로컬 강제 (원격 dispatch 시 ssh 실패) -->
        <key>HEXA_LOCAL</key>
        <string>1</string>
        <key>HEXA_NO_LAUNCHD</key>
        <string>1</string>
    </dict>

    <key>ProcessType</key>
    <string>Background</string>

    <key>Nice</key>
    <integer>15</integer>

    <key>LowPriorityIO</key>
    <true/>
</dict>
</plist>
