# launchd/com.nexus.cross-engine-bus.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.nexus.cross-engine-bus</string>

    <!--
      B5 (2026-04-19) — cross_engine_bus auto-apply launchd.
      Runs cross_engine_bus.hexa --auto-apply in dry-analyze → AUTO-PASS filter → append.
      DISABLED by default; follow-up sessions load it after human review.
    -->

    <key>Disabled</key>
    <true/>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>-c</string>
        <string>/Users/ghost/Dev/nexus/shared/bin/hexa run /Users/ghost/Dev/nexus/shared/n6/cross_engine_bus.hexa --auto-apply</string>
    </array>

    <key>StartInterval</key>
    <integer>86400</integer>

    <key>ThrottleInterval</key>
    <integer>600</integer>

    <key>RunAtLoad</key>
    <false/>

    <key>StandardOutPath</key>
    <string>/Users/ghost/Dev/nexus/shared/logs/cross-engine-bus.stdout.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/ghost/Dev/nexus/shared/logs/cross-engine-bus.stderr.log</string>

    <key>WorkingDirectory</key>
    <string>/Users/ghost/Dev/nexus</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
        <key>HOME</key>
        <string>/Users/ghost</string>
        <key>LANG</key>
        <string>en_US.UTF-8</string>
        <key>NEXUS</key>
        <string>/Users/ghost/Dev/nexus</string>
    </dict>

    <key>ProcessType</key>
    <string>Background</string>

    <key>Nice</key>
    <integer>15</integer>

    <key>LowPriorityIO</key>
    <true/>
</dict>
</plist>
