# launchd/com.nexus6.gen-calc-stubs.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.nexus6.gen-calc-stubs</string>
    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Users/ghost/Dev/nexus6/tools/gen-calc-stubs.sh</string>
    </array>
    <key>WorkingDirectory</key>
    <string>/Users/ghost/Dev/nexus6</string>
    <key>StartInterval</key>
    <integer>900</integer>
    <key>RunAtLoad</key>
    <true/>
    <key>StandardOutPath</key>
    <string>/Users/ghost/Library/Logs/nexus6/gen-calc-stubs.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/ghost/Library/Logs/nexus6/gen-calc-stubs.err</string>
</dict>
</plist>
