# launchd/com.nexus6.dashboard.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key><string>com.nexus6.dashboard</string>
    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Users/ghost/Dev/nexus6/tools/gen-dashboard.sh</string>
    </array>
    <key>WorkingDirectory</key><string>/Users/ghost/Dev/nexus6</string>
    <key>StartInterval</key><integer>600</integer>
    <key>RunAtLoad</key><true/>
    <key>StandardOutPath</key><string>/Users/ghost/Library/Logs/nexus6/dashboard.log</string>
    <key>StandardErrorPath</key><string>/Users/ghost/Library/Logs/nexus6/dashboard.err</string>
</dict>
</plist>
