# launchd/com.nexus.atlas-3d-sync.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.nexus.atlas-3d-sync</string>

    <!--
      ATLAS-P1-4 — atlas.n6 → docs/_mirror_atlas.n6 + _mirror_atlas.live.json bridge.
      H-NOBLOCK 준수: single-pass --watch (변경 없으면 즉시 no-op 종료).
      bootstrap 은 사용자가 명시 요청 시에만:
        launchctl bootstrap gui/$(id -u) shared/launchd/com.nexus.atlas-3d-sync.plist
    -->

    <key>ProgramArguments</key>
    <array>
        <string>/Users/ghost/Dev/nexus/shared/bin/hexa</string>
        <string>run</string>
        <string>/Users/ghost/Dev/nexus/shared/n6/atlas_3d_node_sync.hexa</string>
        <string>--watch</string>
    </array>

    <key>StartInterval</key>
    <integer>30</integer>

    <key>ThrottleInterval</key>
    <integer>10</integer>

    <key>RunAtLoad</key>
    <false/>

    <key>StandardOutPath</key>
    <string>/Users/ghost/Dev/nexus/shared/logs/atlas-3d-sync.stdout.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/ghost/Dev/nexus/shared/logs/atlas-3d-sync.stderr.log</string>

    <key>WorkingDirectory</key>
    <string>/Users/ghost/Dev/nexus</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
        <key>HOME</key>
        <string>/Users/ghost</string>
        <key>LANG</key>
        <string>en_US.UTF-8</string>
        <key>NEXUS</key>
        <string>/Users/ghost/Dev/nexus</string>
    </dict>

    <key>ProcessType</key>
    <string>Background</string>

    <key>Nice</key>
    <integer>15</integer>

    <key>LowPriorityIO</key>
    <true/>
</dict>
</plist>
