# launchd/com.nexus.launcher-cap.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.nexus.launcher-cap</string>
    <key>ProgramArguments</key>
    <array>
        <string>/Users/ghost/.hx/bin/hexa</string>
        <string>/Users/ghost/Dev/airgenome/launcher_cap.hexa</string>
    </array>
    <key>WorkingDirectory</key>
    <string>/Users/ghost/Dev/nexus</string>
    <key>KeepAlive</key>
    <true/>
    <key>ThrottleInterval</key>
    <integer>30</integer>
    <key>StandardOutPath</key>
    <string>/Users/ghost/Library/Logs/nexus/launcher-cap.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/ghost/Library/Logs/nexus/launcher-cap.err</string>
    <key>EnvironmentVariables</key>
    <dict>
        <key>AIRGENOME_CAP_LOG</key>
        <string>/Users/ghost/Library/Logs/nexus/airgenome-cap.jsonl</string>
    </dict>
</dict>
</plist>
