# launchd/com.nexus.atlas-ws.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.nexus.atlas-ws</string>

    <!--
      ATLAS-P3-1 — atlas.n6 → WebSocket live bridge on ws://127.0.0.1:8765.
      DISABLED 2026-04-19: referenced atlas_ws_server.py was removed in R37
      sweep (commit 3c1ef431, python purge policy). KeepAlive caused infinite
      exit=2 respawn. Do NOT bootstrap until hexa-native ws server replaces it.
      Original command (archival only, program path invalid):
        /opt/homebrew/bin/python3 shared/n6/atlas_ws_server.py --atlas shared/n6/atlas.n6
    -->

    <key>Disabled</key>
    <true/>

    <key>ProgramArguments</key>
    <array>
        <string>/opt/homebrew/bin/python3</string>
        <string>/Users/ghost/Dev/nexus/shared/n6/atlas_ws_server.py</string>
        <string>--atlas</string>
        <string>/Users/ghost/Dev/nexus/shared/n6/atlas.n6</string>
    </array>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <dict>
        <key>SuccessfulExit</key>
        <false/>
        <key>Crashed</key>
        <true/>
    </dict>

    <key>ThrottleInterval</key>
    <integer>10</integer>

    <key>StandardOutPath</key>
    <string>/Users/ghost/Dev/nexus/shared/logs/atlas-ws.stdout.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/ghost/Dev/nexus/shared/logs/atlas-ws.stderr.log</string>

    <key>WorkingDirectory</key>
    <string>/Users/ghost/Dev/nexus</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
        <key>HOME</key>
        <string>/Users/ghost</string>
        <key>LANG</key>
        <string>en_US.UTF-8</string>
        <key>NEXUS</key>
        <string>/Users/ghost/Dev/nexus</string>
    </dict>

    <key>ProcessType</key>
    <string>Background</string>

    <key>Nice</key>
    <integer>10</integer>

    <key>LowPriorityIO</key>
    <true/>
</dict>
</plist>
